# supervisord.conf
[supervisord]
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx

[program:uwsgi]
command=/usr/sbin/uwsgi --ini /etc/uwsgi.conf
